# deploy/catvnc-backend.service
[Unit]
Description=CatVNC gateway backend (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=catvnc
Group=catvnc
WorkingDirectory=/opt/catvnc
EnvironmentFile=/opt/catvnc/.env
ExecStart=/opt/catvnc/.venv/bin/uvicorn catvnc.main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=3
LimitNOFILE=65536

# Basic hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/catvnc

[Install]
WantedBy=multi-user.target
